# sh/git-annex/systemd/litelog-sh-git_annex-diskreserve.path
[Unit]
Description=Litelog-Sh Git-Annex Diskreserve Path Monitor
Requires=litelog-sh-systemd-confsync.service

[Path]
PathChanged=/etc/litelog

[Install]
WantedBy=default.target
