# sh/git-annex/systemd/litelog-sh-git_annex-sync.timer
[unit]
Description=Litelog-Sh Git-Annex Sync Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
