# sh/git-annex/systemd/litelog-sh-git_annex-add.path
[Unit]
Description=Litelog-Sh Git-Annex Add Path Monitor
Conflicts=litelog-sh-git_annex-watch.service

[Path]
PathExistsGlob=LOGFILENAME_FINAL * * * * * *
PathChanged=LOGDIR/*

[Install]
WantedBy=default.target
