# sh/git-annex/systemd/litelog-sh-git_annex-diskreserve.path
[Unit]
Description=Litelog-Sh Git-Annex Diskreserve Path Monitor

[Path]
PathChanged=/etc/litelog

[Install]
WantedBy=default.target
