# sh/git/systemd/litelog-sh-git-commit.path
[Unit]
Description=Litelog-Sh Git Commit Path Monitor
Requires=litelog-sh-systemd-confsync.service

[Path]
#PathExistsGlob=LOGFILENAME_FINAL * * * * * *
#PathChanged=LOGDIR/*

[Install]
WantedBy=default.target
