# sh/git/systemd/litelog-sh-git-commit.path
[Unit]
Description=Litelog-Sh Git Commit Path Monitor

[Path]
PathExistsGlob=LOGFILENAME_FINAL * * * * * *
PathChanged=LOGDIR/*

[Install]
WantedBy=default.target
