# sh/gvfs/systemd/litelog-sh-gvfs-move.service
[Unit]
Description=Litelog-Sh GVFS Copying

[Service]
EnvironmentFile=/etc/litelog
Type=oneshot
ExecStart=/bin/sh ${LITELOGDIR}/sh/gvfs/fetchall.sh
