# sh/systemd/litelog-sh-systemd-confsync.path
[Unit]
Description=Litelog-Sh Systemd Configuration Monitor

[Path]
PathChanged=/etc/litelog

[Install]
WantedBy=default.target
