# sh/git/litelog-sh-git-commit.path
[Unit]
Description=Litelog-Sh Git Commit Path Monitor

[Path]
PathChanged=/var/lib/litelog/*

[Install]
WantedBy=default.target
