# sh/git/litelog-sh-git-commit.timer
[Unit]
Description=Litelog-Sh Git Commit Timer

[Timer]
OnCalendar=*:0/2:00:00
Persistent=true

[Install]
WantedBy=default.target
