# sh/git/litelog-sh-git-commit.timer
[Unit]
Description=Litelog-Sh Git Commit Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
