# sh/systemd/systemd/litelog-sh-systemd-confsync.path
[Unit]
Description=Litelog-Sh Systemd Configuration Monitor
Requires=remote-fs.target

[Path]
PathChanged=/etc/litelog

[Install]
WantedBy=default.target
