# sh/video/systemd/litelog-sh-video-compress.path
[Unit]
Description=Litelog-Sh Video Compression Path Monitor

[Path]
PathExistsGlob=LOGFILENAME_INPROGRESS video * * * raw *
PathChanged=LOGDIR/video

[Install]
WantedBy=default.target
