# sh/video/systemd/litelog-sh-video-compress.path
[Unit]
Description=Litelog-Sh Video Compression Path Monitor

[Path]
PathExistsGlob=/var/lib/litelog/video/*_unprocessed.mkv
PathChanged=/var/lib/litelog/video

[Install]
WantedBy=default.target
