# sh/video/systemd/litelog-sh-video-compress.path
[Unit]
Description=Litelog-Sh Video Compression Path Monitor
Requires=litelog-sh-systemd-confsync.service

[Path]
#PathExistsGlob=LOGFILENAME_INPROGRESS video * * * raw *
#PathChanged=LOGDIR/video

[Install]
WantedBy=default.target
