# sh/tight_space/systemd/litelog-sh-tight_space-recording_halt.path
[Unit]
Description=Litelog-Sh Tight Space Recording Halt Path Monitor

[Path]
PathChanged=/var/lib/litelog
PathChanged=/var/lib/litelog/video

[Install]
WantedBy=default.target
