# sh/tight_space/systemd/litelog-sh-tight_space-recording_halt.path
[Unit]
Description=Litelog-Sh Tight Space Recording Halt Path Monitor
Requires=litelog-sh-systemd-confsync.service

[Path]
#PathChanged=LOGDIR
#PathChanged=LOGDIR/*

[Install]
WantedBy=default.target
